# RemindMe/remindme.cron
0  0  *  *  *  /path/to/remindme.cgi send_reminders
